# systemd/system/systemd-vconsole-setup.service.d/override.conf
[Service]
ExecStart=-/usr/bin/setvtrgb /etc/vconsole.rgb
