# udev/rules.d/80-vboxusb.rules
ACTION=="add", SUBSYSTEM=="usb", \
    ENV{DEVTYPE}=="usb_device", \
    RUN+="/usr/bin/find /dev/vboxusb -type d -exec chmod 0755 {} +"
