# systemd/system/local-sdbdown.service
[Unit]
Description=sdbdown
 
[Service]
Type=simple
ExecStart=/home/mirci/bin/sdbdown
 
[Install]
WantedBy=multi-user.target
