# systemd/system/local-lock.service
[Unit]
Description=lock
DefaultDependencies=no
Before=sleep.target

[Service]
Type=forking
ExecStart=/bin/bash -c '/usr/local/bin/lock miro & sleep 2'

[Install]
WantedBy=sleep.target
