# systemd/system/local-hostsblock-fetch.timer
[Unit]
Description=%N timer

[Timer]
OnCalendar=weekly
Persistent=yes

[Install]
WantedBy=timers.target
