# systemd/system/bluetooth.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd --noplugin=sap --compat
#ExecStart=/usr/lib/bluetooth/bluetoothd --noplugin=sap --compat --debug
