# systemd/system/tor.service.d/override.conf
[Service]
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -o tor -g tor -d /var/run/tor
