# systemd/system/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// -l warn -s btrfs --containerd=/run/containerd/containerd.sock
